# ops/systemd/pujia-admin.service
[Unit]
Description=Pujia admin dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=__APP_DIR__
Environment=HOST=127.0.0.1
Environment=PORT=8043
Environment=ADMIN_DATA_PATH=data/redirects.json
EnvironmentFile=__APP_DIR__/.env
ExecStartPre=__APP_DIR__/scripts/sync-basic-auth __APP_DIR__/.env
ExecStart=__APP_DIR__/.venv/bin/waitress-serve --host=${HOST} --port=${PORT} app:app
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
